# roles/docker/files/macvlan_host_bridge.service
[Unit]
After=network.target

[Service]
ExecStart=/sbin/macvlan_host_bridge.sh

[Install]
WantedBy=default.target