# roles/containers/network/swag/files/host-network-bridge.service
[Unit]
After=network.target

[Service]
ExecStart=/sbin/host-network-bridge.sh

[Install]
WantedBy=default.target